# supervisor.ini
[program:autotest_backend]
command=uwsgi --ini /Users/lidongqiang/DjangoProjects/autotest_backend/uwsgi.ini  ;启动程序
priority=1               ;字数越高，优先级越高
numprocs=1               ;启动几个进程
autostart=true           ;随着supervisord启动而启动
autorestart=true         ;程序出错时自动重启
startretries=10          ;启动失败自动重试次数
redirect_stderr=true     ;重定向stderr到stdout
stdout_logfile=/Users/lidongqiang/DjangoProjects/autotest_backend/logs/autotest_backend.log
stderr_logfile=/Users/lidongqiang/DjangoProjects/autotest_backend/logs/autotest_backend.log

[program:autotest_api]
command=gunicorn -c /Users/lidongqiang/DjangoProjects/autotest_api/gunicorn.conf.py autotest_backend.wsgi:application  ;启动程序
priority=1               ;字数越高，优先级越高
numprocs=1               ;启动几个进程
autostart=true           ;随着supervisord启动而启动
autorestart=true         ;程序出错时自动重启
startretries=10          ;启动失败自动重试次数
redirect_stderr=true     ;重定向stderr到stdout
stdout_logfile=/Users/lidongqiang/DjangoProjects/autotest_backend/logs/autotest_backend.log
stderr_logfile=/Users/lidongqiang/DjangoProjects/autotest_backend/logs/autotest_backend.log
